# user/common/files/etc/uci-defaults/zzzz-tailscale-gro
#!/bin/sh
ethtool -K br-lan rx-gro-list on 2>/dev/null
exit 0
